# deploy/stallionussy.service
[Unit]
Description=StallionUSSY - Premium Equine Genetics Exchange
After=network.target

[Service]
Type=simple
User=mojo
Group=mojo
WorkingDirectory=/home/mojo/projects/stallionussy

# The app reads STALLIONUSSY_PORT from the environment when --port is not
# explicitly passed on the command line. Set it in /etc/stallionussy/env to
# override the default (8080). Production typically uses:
#   STALLIONUSSY_PORT=4200
EnvironmentFile=/etc/stallionussy/env
ExecStart=/home/mojo/projects/stallionussy/stallionussy serve

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stallionussy

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/mojo/projects/stallionussy

[Install]
WantedBy=multi-user.target
